# sentry.conf
[program:sentry]
command=uwsgi --ini /etc/sentry.ini
autostart=true
autorestart=true
redirect_stderr=true
